# crontab
*/10 * * * * poetry run python -m from_my_ex
